# dot_config/systemd/user/blueman-applet.service
[Unit]
Description=blueman applet

[Service]
Type=simple
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=xsession.target
