# dot_config/systemd/user/picom.service
[Unit]
Description=picom compositor for X11

[Service]
Type=simple
ExecStart=/usr/bin/picom

[Install]
WantedBy=xsession.target
