# dot_config/systemd/user/redshift.service
[Unit]
Description=redshift compositor for X11

[Service]
Type=simple
ExecStart=/usr/bin/redshift

[Install]
WantedBy=xsession.target
