# dot_config/systemd/user/slack.service
[Unit]
Description=Slack Desktop Client

[Service]
Type=simple
ExecStart=/usr/bin/slack -u

[Install]
WantedBy=xsession.target
