# dot_config/systemd/user/xautolock.service
[Unit]
Description=xautolock

[Service]
Type=simple
ExecStart=/usr/bin/xautolock -locker lock -time 10

[Install]
WantedBy=xsession.target
